# opsworks_agent_monit/templates/default/monit.service.erb
[Unit]
Description=Pro-active monitoring utility for unix systems
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/monit -I
ExecStop=/usr/bin/monit quit
ExecReload=/usr/bin/monit reload
KillMode = process

[Install]
WantedBy=multi-user.target
