# solr-config/templates/default/solr_init.d.sh.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:          solr
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop Solr Server
# Description:       Start/stop Solr Server
### END INIT INFO

SOLR_DIR="/opt/solr-<%= @version %>"
NEWRELIC_PATH="$SOLR_DIR/newrelic/newrelic.jar"
SOLR_BIN="$SOLR_DIR/bin/solr"
START_OPT="<%= @solr_java_mem %><%= " -Dsolr.data.dir=#{@solr_data_dir}  -Duser.timezone=America/Sao_Paulo" unless @solr_data_dir.empty? %>"


# Exit gracefully if the package is not installed
[ -x "$SOLR_BIN" ] || exit 0

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.2-14) to ensure that this file is present
# and status_of_proc is working.
. /lib/lsb/init-functions

#
# Function that starts the daemon/service
#
do_start()
{
  # Return
  #   0 if daemon has been started
  #   1 if daemon was already running
  #   2 if daemon could not be started
  if [ -f $NEWRELIC_PATH ]; then
    NEWRELIC=" -javaagent:$NEWRELIC_PATH "
  fi
  
  SOLR_OPTS=$NEWRELIC $SOLR_BIN start $START_OPT
  exit $!
}

#
# Function that stops the daemon/service
#
do_stop()
{
  $SOLR_BIN stop
}

case "$1" in
  start)
  [ "$VERBOSE" != no ] && log_daemon_msg "Starting"
  do_start
  case "$?" in
    0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
    2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
  esac
  ;;
  stop)
  [ "$VERBOSE" != no ] && log_daemon_msg "Stopping"
  do_stop
  case "$?" in
    0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
    2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
  esac
  ;;
  status)
    $SOLR_BIN status;;

  restart|force-reload)
    log_daemon_msg "Restarting $DESC" "$NAME"
    do_stop
    do_start
  ;;
  *)
  echo "Usage: $0 {start|stop|status|restart|force-reload}" >&2
  exit 3
  ;;
esac

:
