# zabbix_agent/templates/default/zabbix-agent.service.erb
# Do not edit
# Created by chef zabbix-agent::service
[Unit]
Description=zabbix-agent
After=
Requires=

[Service]
TimeoutStartSec=0
ExecStart=/opt/zabbix/sbin/zabbix_agentd

[Install]
WantedBy=
